# deploy/systemd/ai-hunter-producer.service
[Unit]
Description=AI Hunter Queue Producer
After=network.target ai-hunter-api.service
Requires=ai-hunter-api.service

[Service]
Type=simple
WorkingDirectory=/opt/ai-hunter/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/ai-hunter/backend/.venv/bin/python scripts/hunt_queue.py producer --payload-file /opt/ai-hunter/backend/automation_job.json --continuous --enqueue-interval-seconds 60 --max-pending-jobs 3
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
